# eis-install/eis-install.initd
#!/sbin/runscript

depend() {
	need localmount
	after logger
}

start() {
  echo "Starting Eisfair Installation ..."
	/bin/eis-install
	[ -f /reboot ] && reboot
}

stop() {
	echo ""
	exit 0
}
