# eis-install/eis-install.initd
#!/sbin/runscript

depend() {
	need localmount
	after logger
}

start() {
    echo "Starting Eisfair Installation ..."
	/bin/eis-install

}

stop() {
	echo ""
}
