# cuimenu/menus/etc/profile.d/prompt.sh
PS1="\u@\h:\w "
echo "eisfair-ng for alpinelinux"
echo ""
