# cuimenu/menus/etc/profile.d/prompt.sh
PS1="\u@\h:\w "
echo "AlpEisLinux"
echo ""
