# smc-milter-new/smc-milter-new.initd
#!/sbin/runscript
# init.d file for alpine linux.

name=smc-milter-new
daemon=/usr/sbin/$name

depend() {
	need net
	after firewall
}

start() {
    mkdir -p -m0770 /run/milter
    chown -f mail /run/milter
	ebegin "Starting ${name}"
		start-stop-daemon --start --quiet \
			--pidfile /run/${name}.pid \
			--exec ${daemon} -- ${sample_opts}
	eend $?
}

stop() {
	ebegin "Stopping ${name}"
		start-stop-daemon --stop --quiet \
			--pidfile /run/$name.pid \
			--exec ${daemon}
	eend $?
}

