# cui-acme/cui-acme.cron
#!/bin/sh
/usr/bin/acme.sh --cron --home /etc/ssl/acme > /dev/null
