# kbd-mini/kbd-mini.initd
#!/sbin/runscript

depend()
{
	need consolefont
}

start()
{
	[ -f /etc/config.d/base ] && . /etc/config.d/base

	# Set console blank time ESC 9 and VESA powerdown ESC 14
	if [ "0$CONSOLE_BLANK_TIME" -eq 0 ]
	then
		echo -n -e '\033[9;0]\033[14;0]' >/dev/console
	else
		echo -n -e "\033[9;${CONSOLE_BLANK_TIME}]\033[14;${CONSOLE_BLANK_TIME}]" >/dev/console
	fi
}
